# sunshift.service
[Unit]
Description=Sunshift - Automatic color temperature and brightness adjustment
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/andre/.config/hypr/scripts/smooth-temperature.sh

[Install]
WantedBy=default.target 